# rootfs/etc/services.d/ha-integration/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Home Assistant integration service..."

# Wait for main SIP service to be ready
while ! curl -s http://localhost:8088/health > /dev/null 2>&1; do
    bashio::log.info "Waiting for SIP service to be ready..."
    sleep 2
done

# Monitor HA integration health
while true; do
    if ! curl -s http://localhost:8088/health | grep -q '"status":"ok"'; then
        bashio::log.warning "SIP service health check failed"
    fi
    sleep 30
done

